# .config/runit/sv/wl-clip-persist/run
#!/bin/sh

exec 2>&1

[ -r ./conf ] && . ./conf

exec wl-clip-persist --clipboard regular $OPTS 2>&1
