# .config/runit/sv/poweralertd/run
#!/bin/sh

exec 2>&1

[ -r ./conf ] && . ./conf

exec poweralertd -s $OPTS 2>&1
